# packaging/systemd/99-ai-voice-egpu.rules
# Auto-restart the AI Voice Server when an NVIDIA GPU (Vendor 10de) is initialized
ACTION=="add", SUBSYSTEM=="drm", KERNEL=="card*", ATTRS{vendor}=="0x10de", RUN+="/usr/bin/systemctl try-restart --no-block ai-voice-server.service"
